# cloud/oci/ansible/files/qotd.timer
[Unit]
Description=Timer for qotd unit

[Timer]
OnCalendar=*-*-* *:*:10
AccuracySec=1s

[Install]
WantedBy=timers.target
